# assets/systemd/macos-dynamic-wallpaper.service
[Unit]
Description=macOS appearance dinamica (tema+wallpaper chiaro/scuro per ora del giorno)

[Service]
Type=oneshot
ExecStart=@HOME@/.local/bin/macos-appearance.sh auto
