# assets/systemd/macos-dynamic-wallpaper.service
[Unit]
Description=macOS dynamic wallpaper (chiaro/scuro per ora del giorno)

[Service]
Type=oneshot
ExecStart=@HOME@/.local/bin/macos-dynamic-wallpaper.sh
